# packages/mining-fennac/services/fennac-tracker.service
[Unit]
Description=Fennac Share Tracker — dir mining stats
After=fennac.service
Wants=fennac.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/dir/packages/mining-fennac/share-tracker.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
