# packages/mining-dgb-skein/services/dgb-skein-tracker.service
[Unit]
Description=DGB Skein Share Tracker — dir stats
After=dgb-skein.service
Wants=dgb-skein.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/dir/packages/mining-dgb-skein/share-tracker.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
