# packages/mining-dgb-skein/services/dgb-skein-tracker.service
[Unit]
Description=DigiByte Skein Share Tracker — dir
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/dir/packages/mining-dgb-skein/dgb-share-tracker.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
